# misc/systemd/service/org.deepin.linglong.PackageManager.service.in
[Unit]
Description=Linglong dbus service

[Service]
Type=dbus
User=@LINGLONG_USERNAME@
Group=@LINGLONG_USERNAME@
BusName=org.deepin.linglong.PackageManager
ExecStartPre=+@CMAKE_INSTALL_FULL_LIBEXECDIR@/linglong/create-linglong-dirs
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/ll-package-manager
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
