# misc/systemd/system/org.deepin.linglong.SystemHelper.service.in
# SPDX-FileCopyrightText: 2023 UnionTech Software Technology Co., Ltd.
#
# SPDX-License-Identifier: LGPL-3.0-or-later

[Unit]
Description=linglong system helper, only portal can process with root.

[Service]
Type=dbus
User=root
Group=root
BusName=org.deepin.linglong.SystemHelper
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/ll-system-helper
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
